# systemd/bingo-telegram-bot.service
[Unit]
Description=Bingo Telegram Bot (Long Polling)
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/fish_bingo/backend

# Environment variables from .env file
EnvironmentFile=/home/ubuntu/fish_bingo/backend/.env

# Python virtual environment
Environment="PATH=/home/ubuntu/fish_bingo/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/ubuntu/fish_bingo/backend"

# Django settings
Environment="DJANGO_SETTINGS_MODULE=bingo.settings"
Environment="PYTHONUNBUFFERED=1"

# Command to run
ExecStart=/home/ubuntu/fish_bingo/backend/venv/bin/python manage.py runbot

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bingo-telegram-bot

# Security (optional but recommended)
NoNewPrivileges=true
PrivateTmp=true

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target

